# deploy/systemd/crypto-river-24h.timer
[Unit]
Description=Run Crypto Predictor River BTC 24h dashboard forecast

[Timer]
OnCalendar=*-*-* *:06:00
Persistent=true
AccuracySec=5s
Unit=crypto-river-24h.service

[Install]
WantedBy=timers.target
