# deploy/systemd/crypto-river-24h.service
[Unit]
Description=Crypto Predictor River BTC 24h dashboard forecast
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/live
Environment=HOME=/root
Environment=PYTHONUNBUFFERED=1
Environment=BTC_FORECAST_WEB_ROOT=/root/live
Environment=CRYPTO_PREDICTOR_ROOT=/root/claud
ExecStart=/usr/bin/python3 /root/live/crypto_river_24h.py
Nice=10
CPUWeight=50
MemoryMax=1G
TasksMax=64
TimeoutStartSec=3min
NoNewPrivileges=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
