# deploy/systemd/river-1h.timer
[Unit]
Description=Run River hourly shadow forecast after candle close

[Timer]
OnCalendar=*-*-* *:04:00
Persistent=true
AccuracySec=5s
Unit=river-1h.service

[Install]
WantedBy=timers.target
