# deploy/systemd/timesfm-1h.service
[Unit]
Description=TimesFM BTC hourly forecast
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/timesfm-project
Environment=HOME=/root
Environment=PYTHONUNBUFFERED=1
ExecStart=/root/timesfm-env/bin/python /root/timesfm-project/bot_1h.py
NoNewPrivileges=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
Nice=10
CPUWeight=50
MemoryMax=2G
TasksMax=64
TimeoutStartSec=4min
