# deploy/systemd/timesfm-web.service
[Unit]
Description=TimesFM forecast web dashboard
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/live
ExecStart=/usr/bin/python3 /root/live/server.py
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
MemoryMax=128M
TasksMax=32

[Install]
WantedBy=multi-user.target
